# debian/wf-engine.service
[Unit]
Description=Workflow Engine HTTP Server
Documentation=https://github.com/haozicnm/workflow-engine-desktop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/workflow-engine/wf-cli serve --bind 127.0.0.1:19528 --static-dir /opt/workflow-engine/dist
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
PrivateTmp=yes

Environment=BIND=127.0.0.1:19528
Environment=STATIC_DIR=/opt/workflow-engine/dist

[Install]
WantedBy=default.target
